# profile.d/AMDAPPSDK.sh.bak
export AMDAPPSDKROOT="/opt/AMDAPPSDK-3.0"
